# control-plane/radius-stack/supervisor/supervisor-dalocron.conf
[program:supervisor-dalocron]
command=/usr/bin/dalocron.sh
autostart=true
autorestart=false
startretries=0
